# configs/supervisord.conf
[supervisord]
nodaemon=true

[program:stubby]
command=stubby -g

[program:socketserver]
command=nohup /opt/server_dns.py